# packages/deb/after-install.sh
#!/bin/bash
update-rc.d beamium defaults
/etc/init.d/beamium start
