# packages/deb/before-remove.sh
#!/bin/bash
/etc/init.d/beamium stop
update-rc.d -f beamium remove
